# deploy/gce/ib-bot.service
[Unit]
Description=IB Options Bot (US)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=agentbot
Group=agentbot
WorkingDirectory=/opt/agentbot/India
EnvironmentFile=/etc/agentbot/ib.env
ExecStart=/opt/agentbot/venv/bin/python /opt/agentbot/India/ib_options_brain_csv.py
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=45

[Install]
WantedBy=multi-user.target
